# scripts/systemd/clawdbot-auth-monitor.service
[Unit]
Description=ZERO Auth Expiry Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/home/admin/zero/scripts/auth-monitor.sh
# Configure notification channels via environment
Environment=WARN_HOURS=2
# Environment=NOTIFY_PHONE=+1234567890
# Environment=NOTIFY_NTFY=zero-alerts

[Install]
WantedBy=default.target
